# PKGBUILDs/wallpaper-changer/wallpaper-changer.timer
[Unit]
Description=Change wallpaper periodically

[Timer]
OnCalendar=*:00/5:00
AccuracySec=1s
 
[Install]
WantedBy=timers.target
