# PKGBUILDs/eos-brightness/xbacklight.timer
[Unit]
Description=Change brightness hourly

[Timer]
OnStartupSec=10s
OnCalendar=hourly
 
[Install]
WantedBy=timers.target
