# PKGBUILDs/mouse-stuff/mouse-stuff.service
[Unit]
Description=Disable mouse wakeup

[Service]
Type=oneshot
ExecStart=/usr/bin/mouse-stuff
