# PKGBUILDs/mouse-stuff/mouse-stuff.timer
[Unit]
Description=Disable mouse wakeup

[Timer]
OnStartupSec=30sec

[Install]
WantedBy=timers.target
