# PKGBUILDs/m-reminder/m-reminder.timer
[Unit]
Description=Notify something periodically

[Timer]
OnCalendar=hourly
# OnCalendar=*:0/2:0
Persistent=true
 
[Install]
WantedBy=timers.target
